# GRP:archlinux/etc/systemd/system/docker.service.d/imagelocation.conf
[Service]
ExecStart=/usr/bin/dockerd -g /home/docker -H fd://
