# archlinux/etc/systemd/system/docker.service.d/imagelocation.conf
# Copyright 2016 Maximilian Huber <oss@maximilian-huber.de>
# SPDX-License-Identifier: MIT
[Service]
ExecStart=/usr/bin/dockerd -g /home/docker -H fd://
