# GRP:archlinux/etc/udev/rules.d/80-mediasrv-eeti.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/opt/bin/mediaclient --start=5 --systemdcheck"
