# UDOO/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Type=idle
ExecStart=
ExecStart=-/usr/bin/agetty --autologin hubi --noclear %I 38400 linux
